# arc/systemd/gui.service
[Unit]
Description=NFC Attendance Web App
After=network.target

[Service]
Type=simple
User=nakazawa29
WorkingDirectory=/home/nakazawa29/programs/nfc_attendance
EnvironmentFile=/home/nakazawa29/programs/nfc_attendance/.env
ExecStart=/home/nakazawa29/.local/bin/uv run python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
